# debian/DEBIAN/postinst
#!/usr/bin/env sh
set -eu

if ! id nut-counter >/dev/null 2>&1; then
  useradd --system --home /var/lib/nut-counter --shell /usr/sbin/nologin nut-counter
fi

for group_name in gpio video render input; do
  if getent group "$group_name" >/dev/null 2>&1; then
    usermod -aG "$group_name" nut-counter || true
  fi
done

install -d -o nut-counter -g nut-counter /etc/nut-counter /var/lib/nut-counter/models /var/log/nut-counter
install -d -o nut-counter -g nut-counter /var/lib/nut-counter/firefox-profile
install -m 0644 -o nut-counter -g nut-counter /opt/nut-counter/config/firefox-user.js /var/lib/nut-counter/firefox-profile/user.js

# Enterprise policies to force-disable password dialogs and other distractions
install -d /etc/firefox-esr/policies
install -m 0644 /opt/nut-counter/config/firefox-policies.json /etc/firefox-esr/policies/policies.json

chown -R nut-counter:nut-counter /var/lib/nut-counter

if [ -d /etc/sudoers.d ]; then
  echo "nut-counter ALL=(ALL) NOPASSWD: ALL" > /etc/sudoers.d/nut-counter
  chmod 0440 /etc/sudoers.d/nut-counter
fi

if [ ! -f /etc/nut-counter/config.json ]; then
  install -m 0640 -o nut-counter -g nut-counter /opt/nut-counter/config/default-config.json /etc/nut-counter/config.json
else
  chown nut-counter:nut-counter /etc/nut-counter/config.json || true
  chmod 0640 /etc/nut-counter/config.json || true
fi

install -m 0644 /opt/nut-counter/systemd/nut-counter-backend.service /etc/systemd/system/nut-counter-backend.service
install -m 0644 /opt/nut-counter/systemd/nut-counter-kiosk.service /etc/systemd/system/nut-counter-kiosk.service

systemctl daemon-reload
systemctl reenable nut-counter-backend.service || true
systemctl reenable nut-counter-kiosk.service || true
systemctl restart nut-counter-backend.service || true
systemctl restart nut-counter-kiosk.service || true

exit 0
